# scheduler/crontab
* * * * * /scheduler/run-dispatch.sh
0 0 * * * /scheduler/run-generate.sh
